# systemd/ddupdate.timer
[Unit]
Description=Run ddupdate hourly and on boot

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
