# dispatcher.d/50-ddclient
#!/bin/sh

export LC_ALL=C


if [ "$2" = "up" ]; then
  /sbin/ip -o route show dev "$1" | grep -q '^default' && {
     /usr/bin/systemctl start ddclient.service || :
  } || { :; }
fi
