# dispatcher.d/README
Example file for activating ddupdate when an interface goes up. Install into
/etc/NetworkManager/dispatcher.d to activate
